# scripts/tuanzi-photo.service
[Unit]
Description=Tuanzi Photo Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=spi gpio
WorkingDirectory=__WORKING_DIR__
ExecStart=__NODE_BIN__ --max-old-space-size=128 backend/app.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
